# roles/mcp/templates/itential-mcp.service.j2
# Managed by Ansible — itential.mcp role
# Do not edit this file manually. Changes will be overwritten on the next run.

[Unit]
Description=Itential MCP Server
Documentation=https://github.com/itential/itential-mcp
After=network.target

[Service]
Type=simple
User={{ mcp_owner }}
Group={{ mcp_group }}
WorkingDirectory={{ mcp_install_dir }}
EnvironmentFile={{ mcp_env_file }}
ExecStart={{ mcp_venv_dir }}/bin/itential-mcp --config {{ mcp_conf_file }} run
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=itential-mcp

[Install]
WantedBy=multi-user.target
